# crontab
* * * * 5 bash /root/scripts/reset_ss_port.bash